# backup-remote.cron
#!/bin/bash

if [ "$(hostname)" == "Datenteich" ]
then
    /mnt/.bkp/bin/run retry backup-remote /mnt/data.disk/datenteich/ datenteich@192.168.152.11:/mnt/backups.disk/datenteich/ /mnt/backup.disk/datenteich.backups/ bms
fi

if [ "$(hostname)" == "Speichersee" ]
then
    /mnt/.bkp/bin/run retry backup-remote /mnt/data.disk/speichersee/ speichersee@192.168.187.11:/mnt/backups.disk/speichersee/ /mnt/backup.disk/speichersee.backups/ hhs
fi
