# snapshot-create.cron
#!/bin/bash

/mnt/.bkp/bin/run snapshot-create /mnt/data.disk/resilio /mnt/data.disk/resilio.snapshots

if [ "$(hostname)" == "Datenteich" ]
then
    /mnt/.bkp/bin/run snapshot-create /mnt/data.disk/datenteich /mnt/data.disk/datenteich.snapshots
    /mnt/.bkp/bin/run snapshot-create /mnt/boot.disk/rootfs /mnt/boot.disk/rootfs.snapshots
fi

if [ "$(hostname)" == "Speichersee" ]
then
    /mnt/.bkp/bin/run snapshot-create /mnt/data.disk/speichersee /mnt/data.disk/speichersee.snapshots
fi
